# apps/osie-installer-rc
#!/sbin/openrc-run

name="Packet OSIE"
command="/sbin/osie-installer"
command_args="$OSIE_OPTS"
