# apps/runner-rc
#!/sbin/openrc-run

name="Packet Task Runner"
command="/sbin/runner"
