# installer/workflow-helper-rc
#!/sbin/openrc-run

name="Packet Workflow"
command="/sbin/workflow-helper"